# data/eos-config-printer.service.in
[Unit]
Description=Endless Config Printing D-Bus service

[Service]
Type=dbus
BusName=com.endlessm.Config.Printing
ExecStart=@bindir@/eos-config-printer
User=root

